# deluge/supervisord.conf
[supervisord]
nodaemon=true

[program:deluged]
command=/usr/bin/deluged --quiet --do-not-daemonize --config=/config 

[program:deluge-web]
command=/usr/bin/deluge-web --quiet --config=/config
